# .config/systemd/user/mydots-updater.timer
[Unit]
Description=Run mydots updater periodically

[Timer]
OnStartupSec=30 seconds
OnUnitActiveSec=5 minutes

[Install]
WantedBy=timer.target
