# .config/systemd/user/mydots-updater.service
[Unit]
Description=mydots repo updater

[Service]
ExecStart=/usr/bin/git --git-dir %h/mydots --work-tree %h pull --ff-only
