# .config/systemd/user/mydots-updater.timer
[Unit]
Description=Run mydots updater periodically

[Timer]
OnStartupSec=30 seconds
OnUnitInactiveSec=1 hour

[Install]
WantedBy=default.target
